# src/my_project/lib/systemd/system/roscore.service.in
# @ROOTFS@/lib/systemd/system/roscore.service

[Unit]
Description=roscore
After=network.target

[Service]
User=roscore

# https://www.freedesktop.org/software/systemd/man/systemd.exec.html#RuntimeDirectory=
RuntimeDirectory=roscore

# systemd version 229 (shipping with Ubuntu Xenial, as of this writing) doesn't
# export the RUNTIME_DIRECTORY env var, so we have to do this:
Environment=RUNTIME_DIRECTORY=/run/roscore

# systemd version 229 (shipping with Ubuntu Xenial, as of this writing) doesn't
# have any of the following:
#
#StateDirectory=roscore
#CacheDirectory=roscore
#LogsDirectory=roscore
#ConfigurationDirectory=roscore
#
# ...so fake it:

PermissionsStartOnly=true

Environment=STATE_DIRECTORY=/var/lib/roscore
Environment=CACHE_DIRECTORY=/var/cache/roscore
Environment=LOGS_DIRECTORY=/var/log/roscore
Environment=CONFIGURATION_DIRECTORY=/etc/roscore

ExecStartPre=/bin/sh -c 'mkdir -vp ${STATE_DIRECTORY}'
ExecStartPre=/bin/sh -c 'mkdir -vp ${CACHE_DIRECTORY}'
ExecStartPre=/bin/sh -c 'mkdir -vp ${LOGS_DIRECTORY}'
ExecStartPre=/bin/sh -c 'mkdir -vp ${CONFIGURATION_DIRECTORY}'

ExecStartPre=/bin/sh -c 'chown -vR ${USER}:${USER} ${STATE_DIRECTORY}'
ExecStartPre=/bin/sh -c 'chown -vR ${USER}:${USER} ${CACHE_DIRECTORY}'
ExecStartPre=/bin/sh -c 'chown -vR ${USER}:${USER} ${LOGS_DIRECTORY}'
ExecStartPre=/bin/sh -c 'chown -vR ${USER}:${USER} ${CONFIGURATION_DIRECTORY}'

# We need to adapt "The Big Five" env vars into forms that can be easily parsed
# by boost program options; i.e., with a $PROJECT prefix. What's more, we need
# to split the difference between several competing specs that dictate names
# for these things. First read these:
#
# * https://www.gnu.org/prep/standards/html_node/Directory-Variables.html
# * https://specifications.freedesktop.org/basedir-spec/basedir-spec-latest.html
# * https://www.freedesktop.org/software/systemd/man/systemd.exec.html#RuntimeDirectory=
# * https://en.wikipedia.org/wiki/Filesystem_Hierarchy_Standard
#
# Table of Approximal Associations for ficticious Name "foo"
# | systemd                  | xdg                      | gnu
# +----------------------------------------------------------------------------+
# | $RUNTIME_DIRECTORY       | $XDG_RUNTIME_DIR/foo     | $localstatedir/run/foo
# | $STATE_DIRECTORY         | $XDG_CONFIG_HOME/foo     | $localstatedir/lib/foo
# | $CACHE_DIRECTORY         | $XDG_CACHE_HOME/foo      | $localstatedir/cache/foo
# | $LOGS_DIRECTORY          | $XDG_CONFIG_HOME/log/foo | $localstatedir/log/foo
# | $CONFIGURATION_DIRECTORY | $XDG_CONFIG_HOME/foo     | $sysconfdir/foo

# RUNSTATEDIR was added to GNUInstallDirs in cmake v3.9
#
# > RUNSTATEDIR: run-time variable data (LOCALSTATEDIR/run)
#
# -- https://cmake.org/cmake/help/v3.9/module/GNUInstallDirs.html

# This is okay because Ubuntu has symlink /var/run -> /run
Environment=ROSCORE_RUNTIME_DIR=/var/run/roscore
Environment=ROSCORE_STATE_DIR=/var/lib/roscore
Environment=ROSCORE_CACHE_DIR=/var/cache/roscore
Environment=ROSCORE_LOGS_DIR=/var/log/roscore
Environment=ROSCORE_CONFIG_DIR=/etc/roscore

# roscpp logging is... special.
Environment=ROS_HOME=/tmp/roscore
ExecStartPre=/bin/sh -c 'rm -vrf ${ROS_HOME}'

ExecStart=/bin/sh -c '. @CMAKE_INSTALL_PREFIX@/setup.sh && env | sort && roscore -v 2>&1'
ExecStopPost=/bin/sh -c 'rm -vrf ${ROS_HOME}'

# https://www.freedesktop.org/software/systemd/man/systemd-system.conf.html#DefaultStartLimitIntervalSec=
#
# > DefaultStartLimitIntervalSec=, DefaultStartLimitBurst=
# >
# > Configure the default unit start rate limiting, as configured per-service
# > by StartLimitIntervalSec= and StartLimitBurst=. See systemd.service(5) for
# > details on the per-service settings. DefaultStartLimitIntervalSec= defaults
# > to 10s. DefaultStartLimitBurst= defaults to 5.
#
#
# https://www.freedesktop.org/software/systemd/man/systemd.unit.html#StartLimitIntervalSec=interval
#
# > StartLimitIntervalSec=interval, StartLimitBurst=burst
# >
# > Configure unitstart rate limiting. Units which are started more than burst
# > times within an interval time interval are not permitted to start any more.
Restart=always
RestartSec=5

# > core: rename StartLimitInterval= to StartLimitIntervalSec=
#
# https://github.com/systemd/systemd/commit/f0367da7d1a61ad698a55d17b5c28ddce0dc265a#diff-23355de3ac6bfc17cc1269b0de712568
#
# Allow indefinite restarts as long as they don't come at (or faster than) the
# 1/RestartSec frequency as measured over 15 seconds.
StartLimitInterval=14
StartLimitBurst=3

[Install]
WantedBy=default.target
